# deploy/systemd/mine-city-reiki-dictionary.timer
[Unit]
Description=Run Mine City Reiki dictionary maintenance daily between 04:00 and 05:00

[Timer]
OnCalendar=*-*-* 04:00:00
RandomizedDelaySec=3600
Persistent=true

[Install]
WantedBy=timers.target
