# deploy/systemd/mine-city-reiki-sync.timer.example
[Unit]
Description=Run Mine City Reiki sync check hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
